# supervisord.conf
[supervisord]
nodaemon=true

[program:mp4box]
command = MP4Box
autostart = true
